# mailhub/amavisd/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:clamd]
command=/usr/sbin/clamd --foreground
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:amavisd]
command=/usr/sbin/amavisd -c /etc/amavis/amavisd.conf foreground
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# Do not set environment= so amavisd inherits the container env (POSTFIX_HOST, AMAVISD_*, etc. from the stack).
